# install/macos/com.codegroove.secretkeeper.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.codegroove.secretkeeper.agent</string>

    <key>Program</key>
    <string>/Library/PrivilegedHelperTools/secretkeeper-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/PrivilegedHelperTools/secretkeeper-agent</string>
        <string>--config</string>
        <string>/Library/Application Support/SecretKeeper/config.toml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/secretkeeper.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/secretkeeper.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
