# packaging/debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop the service before removal
        if systemctl is-active --quiet secretkeeper.service; then
            echo "Stopping secretkeeper service..."
            systemctl stop secretkeeper.service || true
        fi

        # Disable the service
        if systemctl is-enabled --quiet secretkeeper.service 2>/dev/null; then
            systemctl disable secretkeeper.service || true
        fi
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument: $1" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
