# install/linux/secretkeeper.service
[Unit]
Description=SecretKeeper Secret Exfiltration Prevention Agent
Documentation=https://github.com/secretkeeper/secretkeeper
After=network.target local-fs.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/secretkeeper-agent --config /etc/secretkeeper/config.toml
Restart=on-failure
RestartSec=10s
WatchdogSec=60s

# Run as root (required for CAP_SYS_ADMIN)
User=root
Group=root

# Capabilities
# CAP_SYS_ADMIN: Required for fanotify
# CAP_KILL: Required for SIGSTOP on violation
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_KILL
AmbientCapabilities=CAP_SYS_ADMIN CAP_KILL

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true

# Filesystem protection
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true

# Required writable paths
ReadWritePaths=/var/lib/secretkeeper /var/run

# Memory protection
MemoryDenyWriteExecute=true
LockPersonality=true

# System call filtering
SystemCallFilter=@system-service @mount
SystemCallFilter=~@reboot @swap @obsolete @cpu-emulation @debug @module @raw-io
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=65536
LimitNPROC=64
LimitCORE=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secretkeeper

[Install]
WantedBy=multi-user.target
